# ubuntu/supervisor/extras/conf.d/earlycommands.conf
[program:earlycommands]
priority=2
user=root
startsecs=0
autorestart=false
command=/home/user/earlycommands.sh
