# ubuntu/supervisor/extras/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
user=root
startsecs=0
autorestart=false
priority=1